# deploy/ozone-sfu.service
# Unité systemd pour le nœud média SFU Ozone (processus séparé de l'API).
# Installation :
#   sudo cp target/release/ozone-sfu /usr/local/bin/
#   sudo useradd --system --shell /sbin/nologin ozonesfu
#   sudo cp deploy/ozone-sfu.service /etc/systemd/system/
#   # Fournir le secret partagé (MÊME valeur que l'API) :
#   echo 'OZONE_VOICE_SECRET=<secret-aléatoire-partagé>' | sudo tee /etc/ozone/voice.env
#   sudo systemctl daemon-reload && sudo systemctl enable --now ozone-sfu
#
# OZONE_VOICE_SECRET est REQUIS : sans lui, le SFU refuse toute connexion (fail-closed).

[Unit]
Description=Ozone — nœud média SFU (WebRTC)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ozonesfu
Group=ozonesfu

Environment=OZONE_SFU_BIND=0.0.0.0:8081
Environment=RUST_LOG=info
# Secret partagé avec l'API (obligatoire) — à placer dans ce fichier :
EnvironmentFile=/etc/ozone/voice.env

ExecStart=/usr/local/bin/ozone-sfu
Restart=on-failure
RestartSec=3

# ── Durcissement ──
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
